# bin/cnc.service
[Unit]
Description=CNC
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStartPre=/bin/bash -lc 'until [ -S /tmp/.X11-unix/X0 ]; do sleep 0.2; done'
ExecStartPre=/bin/bash -lc 'until [ -f "$XAUTHORITY" ]; do sleep 0.2; done'
ExecStartPre=/bin/sleep 2
ExecStart=__HOME__/apps/cnc-dro/start.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
